# .crontab
*/5 * * * * cd /var/www/chunta/; bin/chunta slack_status || bin/chunta restart
