# .crontab
PATH=/usr/local/bin:$PATH
0 * * * * cd /var/www/chunta/ && bin/chunta restart
