# .crontab
0 * * * * cd /var/www/chunta/ && bin/chunta restart
